# crontab
*/10 * * * * /bin/sh /srv/scry/snapshot.sh
